# configs/backup/systemd/borgmatic.timer
[Unit]
Description=Run borgmatic backup

[Timer]
OnCalendar=hourly
Persistent=true
RandomizedDelaySec=1min

[Install]
WantedBy=timers.target